# systemd/.config/systemd/user/udiskie.service
[Unit]
Description=Handle automounting of usb devices

[Service]
Type=simple
ExecStart=/usr/bin/udiskie -s

[Install]
WantedBy=default.target
