# systemd/.config/systemd/user/xsession.target
[Unit]
Description=Xession up and running
BindsTo=graphical-session.target
