# systemd-root/etc/systemd/system/ryzenadj.service
[Unit]
Description=ryzenadj Service

[Service]
Type=simple
ExecStart=/usr/bin/ryzenadj --power-saving

[Install]
WantedBy=multi-user.target
